# conf/supervisord.conf
[unix_http_server]
file=/home/work/fcm_app_server/run/supervisor.sock
chmod=0700
chown=work:work

[supervisord]
logfile=/data/logs/fcm_app_server_supervisord.log
logfile_maxbytes=50MB  
logfile_backups=10     
loglevel=info          
pidfile=/home/work/fcm_app_server/run/supervisord.pid 
nodaemon=false
minfds=1024  
minprocs=200 

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/work/fcm_app_server/run/supervisor.sock 

[program:fcm_app_server]
command=/home/work/fcm_app_server/bin/fcm_app_server
process_name=fcm_app_server
numprocs=1
directory=/home/work/fcm_app_server/
redirect_stderr=true
stdout_logfile=/data/logs/fcm_app_server.log
stderr_logfile=/data/logs/fm_app_server.log.err
